# files/etc/cron.d/chevereto
* * * * * www-data php /var/www/html/cli.php -C cron
