# files/etc/cron.d/chevereto
* * * * * www-data php /var/www/html/app/bin/legacy -C cron
